# src/Cattle/cattle_startup.service
[Unit]
After=network.target

[Service]
ExecStart=/usr/local/bin/cattle_startup.sh

[Install]
WantedBy=default.target
