# src/Cattle/startup.service
[Unit]
After=network.target

[Service]
ExecStart=/home/cowmain/project/BarnPi/src/Cattle/startup.sh

[Install]
WantedBy=default.target
